# rootfs/etc/cron.d/wordpress
*/15  *  *  *  * www-data php -f /var/www/wordpress/cron.php
